# wakeup-lid-disable/upstart/wakeup-lid-disable.conf
start on startup
start on runlevel [2345]

pre-start script
    echo LID > /proc/acpi/wakeup
end script

post-stop  script
    echo LID > /proc/acpi/wakeup
end script
